# scripts/ipreg.conf
description "IP address space monitor and registration service"
 
start on runlevel [2345]
stop on runlevel [!2345]


exec start-stop-daemon --start --make-pid --pidfile /var/run/ipreg.pid --background --exec /bin/bash -- -c "/usr/local/bin/ipreg >/var/log/ipreg.log 2>&1"